# docker/nginx/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null

[program:nginx]
command=/usr/local/openresty/bin/openresty -g 'daemon off;'

[program:cron]
command=/usr/sbin/crond -f -d 0
